# addons/codexbar/rootfs/etc/services.d/setup-ui/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "CodexBar setup UI stopped"
