# addons/codexbar/rootfs/etc/services.d/setup-ui/run
#!/usr/bin/with-contenv bashio
set -euo pipefail
export CODEXBAR_HOME=/config
export XDG_CONFIG_HOME=/config/xdg
export CODEXBAR_CONFIG=/config/codexbar/config.json
export HOME=/config
export CODEXBAR_REQUEST_TIMEOUT="$(bashio::config 'request_timeout')"
if bashio::config.has_value 'history_interval'; then
    export CODEXBAR_HISTORY_INTERVAL="$(bashio::config 'history_interval')"
else
    export CODEXBAR_HISTORY_INTERVAL=300
fi
bashio::log.info "Starting CodexBar setup UI on :8099"
exec python3 /usr/share/codexbar-addon/server.py
