# addons/codexbar/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "nginx ingress proxy stopped"
