# addons/codexbar/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
set -euo pipefail
bashio::log.info "Starting nginx ingress proxy on :8099"
exec nginx -g 'daemon off;'
