# addons/codexbar/rootfs/etc/services.d/codexbar/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "CodexBar service stopped"
