# addons/codexbar/rootfs/etc/services.d/codexbar/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

VERSION="$(bashio::config 'codexbar_version')"
REFRESH_INTERVAL="$(bashio::config 'refresh_interval')"
REQUEST_TIMEOUT="$(bashio::config 'request_timeout')"
LOG_LEVEL="$(bashio::config 'log_level')"

export CODEXBAR_HOME=/config
export XDG_CONFIG_HOME=/config/xdg
export CODEXBAR_CONFIG=/config/codexbar/config.json
export HOME=/config

mkdir -p /config/codexbar /config/xdg/codexbar /config/.codex /config/.claude
chmod 700 /config/codexbar /config/.codex /config/.claude || true

# Swift's FileManager.homeDirectoryForCurrentUser resolves /root from passwd
# even when HOME=/config. Keep the official CLI files in persistent /config,
# but expose them at the path CodexBar's Swift credential readers use.
for provider_home in .codex .claude; do
    if [[ -e "/root/${provider_home}" && ! -L "/root/${provider_home}" ]]; then
        rm -rf "/root/${provider_home}"
    fi
    ln -sfn "/config/${provider_home}" "/root/${provider_home}"
done

if [[ ! -x /usr/local/bin/codexbar || ! -f /usr/local/share/codexbar/VERSION || "$(cat /usr/local/share/codexbar/VERSION)" != "${VERSION}" ]]; then
    bashio::log.info "Installing CodexBar CLI ${VERSION} for ${BUILD_ARCH:-unknown arch}"
    /usr/local/bin/install-codexbar-cli "${VERSION}"
fi

bashio::log.info "Starting CodexBar API on 127.0.0.1:8080 (Codex + Claude, refresh=${REFRESH_INTERVAL}s)"
exec /usr/local/bin/codexbar serve \
    --port 8080 \
    --refresh-interval "${REFRESH_INTERVAL}" \
    --request-timeout "${REQUEST_TIMEOUT}" \
    --log-level "${LOG_LEVEL}"
