# addons/codexbar/rootfs/etc/services.d/codexbar/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

VERSION="$(bashio::config 'codexbar_version')"
REFRESH_INTERVAL="$(bashio::config 'refresh_interval')"
REQUEST_TIMEOUT="$(bashio::config 'request_timeout')"
LOG_LEVEL="$(bashio::config 'log_level')"

export CODEXBAR_HOME=/config
export XDG_CONFIG_HOME=/config/xdg
export CODEXBAR_CONFIG=/config/codexbar/config.json
export HOME=/config

mkdir -p /config/codexbar /config/xdg/codexbar /config/.codex /config/.claude
chmod 700 /config/codexbar /config/.codex /config/.claude || true

# Swift's FileManager.homeDirectoryForCurrentUser resolves /root from passwd
# even when HOME=/config. Keep the official CLI files in persistent /config,
# but expose them at the path CodexBar's Swift credential readers use.
root_base="${CODEXBAR_ROOT_HOME:-/root}"
persistent_base="${CODEXBAR_PERSISTENT_HOME:-/config}"
for provider_home in .codex .claude; do
    root_path="${root_base}/${provider_home}"
    persistent_path="${persistent_base}/${provider_home}"
    if [[ -L "${root_path}" && "$(readlink "${root_path}")" == "${persistent_path}" ]]; then
        continue
    fi
    if [[ -e "${root_path}" || -L "${root_path}" ]]; then
        backup_path="${root_path}.pre-codexbar"
        if [[ -e "${backup_path}" || -L "${backup_path}" ]]; then
            bashio::log.fatal "Cannot preserve unexpected ${root_path}; ${backup_path} already exists"
            exit 1
        fi
        bashio::log.warning "Moving unexpected ${root_path} to ${backup_path} before linking persistent storage"
        mv "${root_path}" "${backup_path}"
    fi
    ln -s "${persistent_path}" "${root_path}"
done

if [[ ! -x /usr/local/bin/codexbar || ! -f /usr/local/share/codexbar/VERSION || "$(cat /usr/local/share/codexbar/VERSION)" != "${VERSION}" ]]; then
    bashio::log.info "Installing CodexBar CLI ${VERSION} for ${BUILD_ARCH:-unknown arch}"
    /usr/local/bin/install-codexbar-cli "${VERSION}"
fi

bashio::log.info "Starting CodexBar API on 127.0.0.1:8080 (Codex + Claude, refresh=${REFRESH_INTERVAL}s)"
exec /usr/local/bin/codexbar serve \
    --port 8080 \
    --refresh-interval "${REFRESH_INTERVAL}" \
    --request-timeout "${REQUEST_TIMEOUT}" \
    --log-level "${LOG_LEVEL}"
