# addons/meshcore-migration-ota/rootfs/etc/services.d/meshcore-ota-worker/finish
#!/usr/bin/execlineb -S0
s6-svscanctl -t /run/service
