# addons/meshcore-migration-ota/rootfs/etc/services.d/meshcore-ota-worker/run
#!/usr/bin/with-contenv bashio
set -euo pipefail
umask 077
unset PYTHONPATH PYTHONHOME PYTHONSTARTUP PYTHONINSPECT
export PYTHONDONTWRITEBYTECODE=1

exec /usr/local/bin/meshcore-ota-worker
